# config/tether.service
[Unit]
Description=Tether Server v3 — 极简消息中转
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=zzsky
ExecStart=/home/zzsky/.hermes/tether/tether_server.py --port 9001
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
