# config/tether-watcher.service
[Unit]
Description=Tether Watcher — 2秒轮询处理消息
After=tether.service
# 注意：不使用 BindsTo，Tether Server 和 Watcher 是两个独立进程，
# 任一重启不应连带杀死对方。启动顺序由 After 保证。

[Service]
Type=simple
ExecStart=/home/zzsky/.hermes/tether/tether_watcher.py
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
